# services/avahi/run
#!/bin/bash

exec avahi-daemon --no-chroot
