# services/avahi/run
#!/bin/bash
echo "Starting Avahi daemon"
exec avahi-daemon --no-chroot
