# services/plex/run
#!/bin/bash

exec /sbin/setuser abc /usr/sbin/start_pms
