# services/plex/run
#!/bin/bash

exec /usr/bin/chpst -u abc /usr/sbin/start_pms
