# services/plex/run
#!/bin/bash

echo "Starting Plex Media Server."
exec /sbin/setuser abc /usr/sbin/start_pms
