# services/plex/run
#!/bin/bash

exec /sbin/setuser plex /usr/sbin/start_pms
