# root/etc/services.d/dbus/run
#!/usr/bin/with-contenv bash

echo "Starting dbus-daemon"
exec dbus-daemon --system --nofork
