# root/etc/services.d/plex/run
#!/usr/bin/with-contenv bash

echo "Starting Plex Media Server."
exec \
	s6-setuidgid abc /bin/bash -c \
	'LD_LIBRARY_PATH=/usr/lib/plexmediaserver /usr/lib/plexmediaserver/Plex\ Media\ Server'
