# root/etc/services.d/plex/run
#!/usr/bin/with-contenv bash

echo "Starting Plex Media Server."
exec \
	s6-setuidgid abc /bin/sh -c \
	'/usr/lib/plexmediaserver/Plex\ Media\ Server'


